# extra-coreos/usr/lib/dracut/modules.d/40ignition-ostree/boot-part-discovery.service
[Unit]
Description=Discover Boot Partition Availability
ConditionPathExists=/etc/initrd-release
Before=boot-part-discovered.target
After=systemd-udev-settle.service
Wants=systemd-udev-settle.service
Before=sysroot.mount

# After=local-fs-pre.target
# Before=local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/libexec/boot-part-discovery.sh
RemainAfterExit=yes
